# server/deploy/job-dashboard-api.service.example
[Unit]
Description=Job Dashboard API Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Set the runtime user/group for your VM (example: ubuntu)
User=ubuntu
Group=ubuntu

# Set your checked-out server directory
WorkingDirectory=/home/ubuntu/dashboard/server

# Keep production secrets and runtime config outside the repo
EnvironmentFile=/etc/job-dashboard/api.env

# If Node is installed with nvm, use the full absolute node path
# Example below matches Oracle VM setup with nvm-managed Node 20.
ExecStart=/home/ubuntu/.nvm/versions/node/v20.20.1/bin/node index.js
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

# Because code lives under /home/ubuntu, true can block service access.
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
